# files/etc/cron.d/website_cron
* * * * * root curl -b /tmp/cookie.txt -c /tmp/cookie.txt -sS ${GK_SITE_BASE_SERVER_URL}/cron
