# files/etc/cron.d/website_cron
* * * * * root curl -b /tmp/cookie.txt -c /tmp/cookie.txt -H "User-Agent: GeoKretyCron/1.0.1" -sS ${GK_SITE_BASE_SERVER_URL}/cron
